# cron/crontab
0 0 * * * (/scripts/command.sh 2>&1) >/var/log/cron.log
